# supervisor/rapidbistro-worker.conf
[program:rapidbistro-worker]
; ── RapidBistro — Queue Worker (Supervisor) ──────────────────────────────────
; Copia este archivo a /etc/supervisor/conf.d/rapidbistro-worker.conf
; Luego ejecuta: sudo supervisorctl reread && sudo supervisorctl update
; Para ver estado: sudo supervisorctl status

command=php /var/www/rapidbistro-app/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
directory=/var/www/rapidbistro-app
user=www-data
numprocs=2
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=300
stopasgroup=true
killasgroup=true

; Logs
stdout_logfile=/var/log/supervisor/rapidbistro-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/rapidbistro-worker-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

; Restart workers diariamente para evitar memory leaks
; El flag --max-time=3600 hace que el worker se detenga solo después de 1h
; y Supervisor lo reinicia automáticamente.
