# supervisor/rapidbistro-scheduler.conf
[program:rapidbistro-scheduler]
; ── RapidBistro — Laravel Scheduler (Supervisor) ─────────────────────────────
; Alternativa a cron para entornos donde no se quiere tocar el crontab del servidor.
; En producción es más común usar crontab, ver deploy.md.
; Este programa ejecuta el scheduler cada minuto de forma continua.

command=/bin/bash -c "while true; do php /var/www/rapidbistro-app/artisan schedule:run >> /dev/null 2>&1; sleep 60; done"
directory=/var/www/rapidbistro-app
user=www-data
autostart=true
autorestart=true
startsecs=0
startretries=3

stdout_logfile=/var/log/supervisor/rapidbistro-scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/rapidbistro-scheduler-error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
